# ansible/roles/jellyfin/files/ts-jellyfin.container
[Unit]
Description=Tailscale Container for Jellyfin

[Container]
Pod=jellyfin.pod
Image=docker.io/tailscale/tailscale:latest
AutoUpdate=registry
ContainerName=ts-jellyfin
HostName=jellyfin
Secret=ts-authkey,type=env,target=TS_AUTHKEY
Environment=TS_STATE_DIR=/var/lib/tailscale
Environment=TS_USERSPACE=false
Environment=TS_SERVE_CONFIG=/config/serve-config.json
Volume=ts-jellyfin-state.volume:/var/lib/tailscale:Z
Volume=./serve-config.json:/config/serve-config.json:Z,ro
AddDevice=/dev/net/tun:/dev/net/tun
AddCapability=NET_ADMIN
